# cron_example/create_home.cron
53 * * * * /home/XXX/Dracoon-Automation-Scripts/cronjob_scripts/create_home_folder.sh
